# contrib/init/gleecbtcd.conf
description "GleecBTC Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env GLEECGBCD_BIN="/usr/bin/gleecbtcd"
env GLEECGBCD_USER="gleecbtc"
env GLEECGBCD_GROUP="gleecbtc"
env GLEECGBCD_PIDDIR="/var/run/gleecbtcd"
# upstart can't handle variables constructed with other variables
env GLEECGBCD_PIDFILE="/var/run/gleecbtcd/gleecbtcd.pid"
env GLEECGBCD_CONFIGFILE="/etc/gleecbtc/gleecbtc.conf"
env GLEECGBCD_DATADIR="/var/lib/gleecbtcd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # gleecbtcd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$GLEECGBCD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run gleecbtcd."
        echo "The setting must appear in $GLEECGBCD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"GleecBTC Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$GLEECGBCD_PIDDIR"
    chmod 0755 "$GLEECGBCD_PIDDIR"
    chown $GLEECGBCD_USER:$GLEECGBCD_GROUP "$GLEECGBCD_PIDDIR"
    chown $GLEECGBCD_USER:$GLEECGBCD_GROUP "$GLEECGBCD_CONFIGFILE"
    chmod 0660 "$GLEECGBCD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$GLEECGBCD_PIDFILE" \
    --chuid $GLEECGBCD_USER:$GLEECGBCD_GROUP \
    --exec "$GLEECGBCD_BIN" \
    -- \
    -pid="$GLEECGBCD_PIDFILE" \
    -conf="$GLEECGBCD_CONFIGFILE" \
    -datadir="$GLEECGBCD_DATADIR" \
    -disablewallet \
    -daemon

